# apps/hetzner/wireguard/files/etc/systemd/system/wg-quick-watcher@.path
[Unit]
Description=Watcher for WireGuard via wg-quick(8) for %I

[Path]
PathChanged=/etc/wireguard/%i.conf

[Install]
WantedBy=multi-user.target
